# sync/sync.timer
[Unit]
Description=Sync daily

[Timer]
OnCalendar=*-*-* 00:00:01
Unit=sync.service
Persistent=true

[Install]
WantedBy=timers.target
